# deploy/com.ianduclos.gridmapper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.ianduclos.gridmapper</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>/Users/ianduclos/_SecondBrain/01_Projects/gridmapper/node_modules/tsx/dist/cli.mjs</string>
		<string>/Users/ianduclos/_SecondBrain/01_Projects/gridmapper/src/cli/sim.ts</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/ianduclos/_SecondBrain/01_Projects/gridmapper</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/Users/ianduclos/Library/Logs/gridmapper.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/ianduclos/Library/Logs/gridmapper.log</string>
</dict>
</plist>
